# deploy/supervisor/patronr-messenger-payments.conf
[program:patronr-messenger-payments]
directory=/path/to/patronr
command=php bin/console messenger:consume payments --time-limit=3600 --memory-limit=256M --no-reset
user=your-user
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stdout_logfile=/path/to/patronr/var/log/supervisor-payments.log
stderr_logfile=/path/to/patronr/var/log/supervisor-payments-error.log
environment=APP_ENV="prod"
