# deploy/supervisor/patronr-messenger-integrations.conf
[program:patronr-messenger-integrations]
directory=/path/to/patronr
command=php bin/console messenger:consume integrations --time-limit=3600 --memory-limit=256M --no-reset
user=your-user
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stdout_logfile=/path/to/patronr/var/log/supervisor-integrations.log
stderr_logfile=/path/to/patronr/var/log/supervisor-integrations-error.log
environment=APP_ENV="prod"
